# docs/server-configs/supervisor.conf
; =============================================================================
; Supervisor Configuration — CBT SMK LMS
; Path: /etc/supervisor/conf.d/smk-lms.conf
; Apply: sudo supervisorctl reread && sudo supervisorctl update
; =============================================================================

; === Default Queue Workers (3 procs) ===
; Handles: grading jobs, notifications, exports, general tasks
[program:laravel-worker-default]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/smk-lms/artisan queue:work redis --queue=default --sleep=3 --tries=3 --max-time=3600 --memory=128 --backoff=10,30,60
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
numprocs=3
user=www-data
redirect_stderr=true
stdout_logfile=/var/log/supervisor/worker-default.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stopwaitsecs=30

; === Exam Persist Workers (2 procs) ===
; Dedicated for PersistAnswersJob — separated so general jobs don't block persistence
[program:laravel-worker-exam]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/smk-lms/artisan queue:work redis --queue=exam-persist --sleep=3 --tries=5 --max-time=3600 --memory=128 --backoff=5,15,30
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
numprocs=2
user=www-data
redirect_stderr=true
stdout_logfile=/var/log/supervisor/worker-exam.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stopwaitsecs=30

; === Laravel Reverb (WebSocket Server) ===
[program:reverb]
process_name=%(program_name)s
command=php /var/www/smk-lms/artisan reverb:start --host=127.0.0.1 --port=8080
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
numprocs=1
user=www-data
redirect_stderr=true
stdout_logfile=/var/log/supervisor/reverb.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stopwaitsecs=10

; === Laravel Scheduler (cron replacement via Supervisor) ===
[program:laravel-scheduler]
process_name=%(program_name)s
command=/bin/bash -c "while true; do php /var/www/smk-lms/artisan schedule:run --no-interaction >> /dev/null 2>&1; sleep 60; done"
autostart=true
autorestart=true
numprocs=1
user=www-data
redirect_stderr=true
stdout_logfile=/var/log/supervisor/scheduler.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=3

; === Group: All SMK LMS processes ===
[group:smk-lms]
programs=laravel-worker-default,laravel-worker-exam,reverb,laravel-scheduler
priority=999
